# borsch-shop/docs/server/landing.service
[Unit]
Description=Borsch Landing (Next.js Standalone)
After=network.target

[Service]
Type=simple
User=root
# Путь на сервере к приложению. Меняем на тот, где будет лежать код.
WorkingDirectory=/var/www/borsch-shop/landing/.next/standalone
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
