# RMS/docs/server/pocketbase.service
[Unit]
Description=PocketBase Server (RMS AI OS)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/rms-ai-os/backend
# Предполагается, что бинарь pocketbase лежит в /var/www/rms-ai-os/backend
ExecStart=/var/www/rms-ai-os/backend/pocketbase serve --http=127.0.0.1:8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
